# smf/wzfsadm.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle
  SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    Manifest created by svcbundle (2016-Apr-13 13:35:23+0000)
-->
<service_bundle type="manifest" name="system/filesystem/wzfsadm">
    <service version="1" type="service" name="system/filesystem/wzfsadm">
        <!--
            The following dependency keeps us from starting until the
            multi-user milestone is reached.
        -->
        <dependency restart_on="none" type="service"
            name="multi_user_dependency" grouping="require_all">
            <service_fmri value="svc:/milestone/multi-user"/>
        </dependency>

	<method_context>
	    <method_credential user='wzfsadm' group='staff' />
	</method_context>

        <exec_method 
		timeout_seconds="60"
		type="method"
		name="start"
		exec="/lib/svc/method/wzfsadm start"/>
        <!--
            The exec attribute below can be changed to a command that SMF
            should execute to stop the service.  See smf_method(5) for more
            details.
        -->
        <exec_method 
		timeout_seconds="60"
		type="method"
		name="stop"
		exec="/lib/svc/method/wzfsadm stop"/>
        <!--
            The exec attribute below can be changed to a command that SMF
            should execute when the service is refreshed.  Services are
            typically refreshed when their properties are changed in the
            SMF repository.  See smf_method(5) for more details.  It is
            common to retain the value of :true which means that SMF will
            take no action when the service is refreshed.  Alternatively,
            you may wish to provide a method to reread the SMF repository
            and act on any configuration changes.
        -->     
	
        <!--
            We do not need a duration property group, because contract is
            the default.  Search for duration in svc.startd(1M).
        -->
        <instance enabled="true" name="default"/>        
    </service>
</service_bundle>
